# imageroot/systemd/user/openvpnui.service
#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

#
# This systemd unit starts a openvpnui instance using Podman.
# Most parts of this file come from podman-generate-systemd.
#

[Unit]
Description=Podman openvpnui.service
Requires=openvpnui-app.service
Before=openvpnui-app.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-%S/state/environment
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/openvpnui.pid %t/openvpnui.pod-id
ExecStartPre=/usr/bin/podman pod create --infra-conmon-pidfile %t/openvpnui.pid \
    --pod-id-file %t/openvpnui.pod-id \
    --name openvpnui \
    --publish 127.0.0.1:${TCP_PORT}:80 \
    --replace
ExecStart=/usr/bin/podman pod start --pod-id-file %t/openvpnui.pod-id
ExecStop=/usr/bin/podman pod stop --ignore --pod-id-file %t/openvpnui.pod-id -t 10
ExecStopPost=/usr/bin/podman pod rm --ignore -f --pod-id-file %t/openvpnui.pod-id
PIDFile=%t/openvpnui.pid
Type=forking

[Install]
WantedBy=default.target
